# examples/macos/com.streetkings.dismediapresence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.streetkings.dismediapresence</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/DisMediaPresence/.venv/bin/dmp</string>
    <string>--settings</string>
    <string>/ABSOLUTE/PATH/TO/DisMediaPresence/dmp.settings</string>
    <string>run</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/DisMediaPresence</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/DisMediaPresence/logs/dmp.out.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/DisMediaPresence/logs/dmp.err.log</string>
</dict>
</plist>
